# com.mikesimonsen.fed-expectations.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mikesimonsen.fed-expectations</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mikesimonsen/projects/econ-data/fed_expectations.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mikesimonsen/projects/econ-data</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
